# attendance_scheduler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd job: daily attendance CSV download at 9:00 AM

  DO NOT edit this file directly with paths — run setup_launchd.sh instead.
  That script fills in __PROJECT_DIR__ and __VENV_PYTHON__ automatically.

  To install manually:
    1. Replace __PROJECT_DIR__ with the absolute path to attendance_automation/
    2. Replace __VENV_PYTHON__ with the absolute path to venv/bin/python
    3. cp attendance_scheduler.plist ~/Library/LaunchAgents/com.attendance.automation.plist
    4. launchctl load ~/Library/LaunchAgents/com.attendance.automation.plist
-->
<plist version="1.0">
<dict>

    <!-- Unique identifier for this job -->
    <key>Label</key>
    <string>com.attendance.automation</string>

    <!-- Command to run: venv Python + script path -->
    <key>ProgramArguments</key>
    <array>
        <string>__VENV_PYTHON__</string>
        <string>__PROJECT_DIR__/main.py</string>
    </array>

    <!-- Run daily at 18:00 UK local time (Mac must be set to London / Europe timezone) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>18</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Working directory so that .env is found relative to the script -->
    <key>WorkingDirectory</key>
    <string>__PROJECT_DIR__</string>

    <!-- Capture output for troubleshooting -->
    <key>StandardOutPath</key>
    <string>__LOG_DIR__/launchd_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__LOG_DIR__/launchd_stderr.log</string>

    <!-- Do NOT run immediately on load — only at the scheduled time -->
    <key>RunAtLoad</key>
    <false/>

    <!-- launchd provides a stripped environment; restore minimum PATH and HOME -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>__HOME__</string>
    </dict>

</dict>
</plist>
